# net-proxy/goagent/files/goagent-initd
#!/sbin/runscript
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

extra_started_commands="update"

GOAGENT_PID="--pidfile /var/run/goagent.pid"
GOAGENT_LOG="/var/log/goagent.log"

depend() {
    use logger
}

start() {
    ebegin "Starting GoAgent"
    start-stop-daemon --start --name goagent \
    $GOAGENT_PID --make-pid \
    --exec python3 /opt/goagent/proxy.py &> $GOAGENT_LOG &
    eend $?
}

stop() {
    ebegin "Stopping GoAgent"
    start-stop-daemon --stop --name goagent $GOAGENT_PID \
    --exec python3 /opt/goagent/proxy.py
    eend $?
}

